# services/nscd/run
#!/usr/bin/with-contenv sh
echo >&2 "Starting unscd..."
exec /usr/sbin/nscd -f /etc/nscd.conf
