# services/nscd/run
#!/usr/bin/with-contenv sh
echo >&2 "Starting unscd..."
exec sudo -u unscd /usr/sbin/nscd -d -f /etc/nscd.conf
